# root/lib/systemd/system/tbbuttonsd.service
[Unit]
Description=tbbuttonsd
After=network.target

[Service]
ExecStart=/usr/local/bin/tbbuttonsd

[Install]
WantedBy=multi-user.target
