# root/lib/systemd/system/tbprocessd.service
[Unit]
Description=tbprocessd

[Service]
Environment=HOME_APP=/apps/home
ExecStart=/usr/local/bin/tbprocessd

[Install]
WantedBy=multi-user.target
