# root/DEBIAN/prerm
#!/bin/sh

set -o errexit

if [ "$1" = "remove" -o "$1" = "abort-install" -o "$1" = "disappear" ]; then
    rm -f /boot/overlays/tingbot-overlay.dtb
    update-rc.d -f tbwifisetup remove
fi