# root/etc/init.d/tbprocessd.service
[Unit]
Description=tbprocessd

[Service]
ExecStart=/usr/local/bin/tbprocessd

[Install]
WantedBy=multi-user.target
